# scripts/systemd/bumblebeed.service.in
[Unit]
Description=Bumblebee C Daemon

[Service]
Type=simple
ExecStartPre=-/usr/bin/prime-select intel
ExecStart=@SBINDIR@/bumblebeed --use-syslog
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
