# scripts/systemd/bumblebeed_wild.service.in
[Unit]
Description=Bumblebee C Daemon

[Service]
Type=simple
ExecStartPre=-/usr/bin/prime-select intel
ExecStart=@SBINDIR@/bumblebee.bash @BBCONFDIR@/bumblebee.conf
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
