# scripts/systemd/bumblebeed.service.in
[Unit]
Description=Bumblebee C Daemon

[Service]
Type=simple
ExecStart=@SBINDIR@/bumblebeed
Restart=always
RestartSec=60
StandardOutput=kmsg

[Install]
WantedBy=multi-user.target
