# scripts/upstart/bumblebeed.conf.in
# Upstart job for Bumblebee daemon

description "Bumblebee supporting nVidia Optimus cards"
author      "Lekensteyn <lekensteyn@gmail.com>"

# XXX: do we need to stop before / after a graphical DE like kdm?
start on    (runlevel [2345])
stop on     (runlevel [016])

# restart in case it crashed for some reason
respawn

expect fork
exec '@BINDIR@/bumblebeed' --daemon
